# etc/supervisor.d/ancol-go-worker.ini
[group:ancol]
programs=horizon
priority=25

[program:horizon]
command=php /app/artisan horizon
process_name=%(program_name)s_%(process_num)02d
numprocs=1
priority=999
startsecs=0
startretries=3
autostart=true
autorestart=true
user=ipnet
redirect_stderr=true
stdout_logfile=/app/storage/logs/horizon.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/storage/logs/horizon-err.log
stderr_logfile_maxbytes=0
